# etc/jabberd2-router.service.in
[Unit]
Description=Jabber Router XML Packet Distributior
After=network.target
BindsTo=@package@.service

[Service]
User=jabber
ExecStart=@bindir@/@package_router_bin@ -c @configdir@/router.xml
Restart=on-failure

[Install]
WantedBy=multi-user.target
