# etc/jabberd2-sm.service.in
[Unit]
Description=Jabber IM Session Manager
After=network.target
BindsTo=@package@.service

[Service]
User=jabber
ExecStart=@bindir@/@package_sm_bin@ -c @configdir@/sm.xml
Restart=on-failure

[Install]
WantedBy=multi-user.target
