# etc/jabberd2-sm.service.in
[Unit]
Description=Jabber IM Session Manager
After=network.target
BindTo=@package@.service

[Service]
User=jabber
ExecStart=@bindir@/@package_sm_bin@ -c @configdir@/sm.xml

[Install]
WantedBy=multi-user.target
